# roles/2cfg.lb.configure_nginx_server/templates/rc.local
#!/bin/bash

ls -1 /etc/rc.local.d/* | bash

exit 0;